# ops/monitoring/minhoo-idempotency-alert.service
[Unit]
Description=Minhoo Idempotency Monitor Failure Alert
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=appuser
Group=appuser
WorkingDirectory=/var/www/minhoo-api/minhoo_api
Environment=NODE_ENV=production
Environment="FEED_CACHE_ALERT_TITLE=[Minhoo] Idempotency monitor FAILED"
ExecStart=/usr/bin/node /var/www/minhoo-api/minhoo_api/scripts/notify-feed-cache-alert.js --unit minhoo-idempotency-monitor.service --since-minutes 30 --lines 160
