# ops/monitoring/minhoo-risk-monitor.timer
[Unit]
Description=Run Minhoo Risk Monitor every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=minhoo-risk-monitor.service

[Install]
WantedBy=timers.target

