# ops/monitoring/minhoo-feed-slo-sampler.service
[Unit]
Description=Minhoo Feed SLO Sampler (every 5 min)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=appuser
Group=appuser
WorkingDirectory=/var/www/minhoo-api/minhoo_api
Environment=NODE_ENV=production
Environment=FEED_SLO_STRICT=1
ExecStart=/usr/bin/node /var/www/minhoo-api/minhoo_api/scripts/monitor-feed-slo-sample.js

