# ops/prod-scale/minhoo-api.service
[Unit]
Description=Minhoo API (Node.js) - Blue
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/var/www/minhoo-api/minhoo_api
EnvironmentFile=/var/www/minhoo-api/minhoo_api/.env
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /var/www/minhoo-api/minhoo_api/dist/app.js
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
